# shanoir-ng-template/DockerWithJdk8AndMariaDb/supervisor/mysql.conf
[program:mariadb]
command=run-mysql
stderr_logfile = /vol/log/supervisor/mariadb-stderr.log
stdout_logfile = /vol/log/supervisor/mariadb-stdout.log

