# shanoir-ng-template/shanoir-ng-base-ms/files/etc/supervisor/conf.d/msapp.conf
[program:msapp]
command=run-msapp
stderr_logfile = /vol/log/supervisor/msapp-stderr.log
stdout_logfile = /vol/log/supervisor/msapp-stdout.log
