# shanoir-ng-users/src/main/docker/ms.users.sv.dev.conf
[program:usersmsapp]
command=run-shanoir-ng-users
stderr_logfile = /vol/log/supervisor/usersmsapp-stderr.log
stdout_logfile = /vol/log/supervisor/usersmsapp-stdout.log
